# debian/torouterui.postrm
if [ "$1" = "purge" ] ; then
	update-rc.d torouterui remove >/dev/null
fi
